# monitor/fuzzel-monitor-scheme.path
[Unit]
Description=Sync fuzzel and caelestia schemes

[Path]
PathModified=%C/caelestia/scheme/current.txt
Unit=fuzzel-monitor-scheme.service

[Install]
WantedBy=default.target
